# debian/postinst
# WebUI key for the API
/usr/local/skysql/config/generateAPIkey.sh 1
/usr/local/skysql/config/generateAPIkey_json.sh 1  /usr/local/skysql/
# Monitor key for the API
/usr/local/skysql/config/generateAPIkey.sh 3
# Restart the Monitor so that it reads the new key
/etc/init.d/mariadb-enterprise-monitor restart

# setup httpd start and restart httpd
#sed -i 's/# chkconfig: -/# chkconfig: 2345/' /etc/init.d/httpd
#rm -f /etc/rc{2,3,4,5}.d/K*httpd*
update-rc.d apache2 defaults
/etc/init.d/apache2 restart

# Cleanup
rm -f  /usr/local/skysql/config/generateAPIkey*.sh
